# worker/supervisord.conf
[supervisord]
nodaemon=true

[program:slurmd]
command=/usr/local/sbin/slurmd
autostart=true
autorestart=true
